# root/etc/s6-overlay/s6-rc.d/named/run
#!/usr/bin/with-contenv bash

exec /usr/sbin/named -g -u bind -c /config/named.conf
